# wotd
#!/sbin/openrc-run

depend() {
    need net
}

command="/root/word_of_the_day_endpoint/python3 app.py"
pidfile="/run/${RC_SVCNAME}/pid"
