# wotd
#!/sbin/openrc-run

depend() {
    need net
}

command="/usr/bin/python3 /root/word_of_the_day_endpoint/app.py &"
pidfile="/run/${RC_SVCNAME}/pid"
